# deploy/upload-lift-commission.service
[Unit]
Description=Upload Studio - Commission Auto-Charge Worker
After=network.target upload-studio.service
Requires=upload-studio.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/upload-studio
ExecStart=/usr/bin/npx tsx workers/commission.worker.ts
Restart=always
RestartSec=30
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=upload-studio-commission

[Install]
WantedBy=multi-user.target
