# deploy/upload-lift-preflight.service
[Unit]
Description=Upload Studio - Preflight Worker
After=network.target redis.service postgresql.service upload-studio.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/upload-studio
ExecStart=/usr/bin/pnpm run worker:preflight
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/var/www/upload-studio/.env

# Logging
StandardOutput=append:/var/log/upload-studio/preflight-worker.log
StandardError=append:/var/log/upload-studio/preflight-worker-error.log

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/www/upload-studio/temp /var/log/upload-studio

[Install]
WantedBy=multi-user.target

